# ansible/roles/sentry/templates/supervisor.conf.j2
[program:sentry]
command={{ app_home }}/start_sentry.sh
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile = {{ app_log }}/supervisor.log 
