# ansible/roles/sentry/templates/uwsgi.upstart
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

env LOGTO=/var/log/uwsgi/emperor.log

exec uwsgi  --die-on-term --logto $LOGTO --emperor /etc/uwsgi/apps-enabled